# ansible/roles/deploy/templates/eveleve.service.j2
[Unit]
Description=EvelEve Discord Bot

[Service]
Type=simple
ExecStart=/usr/bin/eveleve master
Restart=always

[Install]
WantedBy=multi-user.target

